# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery status notifications

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-notify
